# deploy/systemd/terminal-bot.service.d/override.conf
# Скопировать на сервер:
#   sudo mkdir -p /etc/systemd/system/terminal-bot.service.d
#   sudo cp deploy/systemd/terminal-bot.service.d/override.conf /etc/systemd/system/terminal-bot.service.d/
#   sudo systemctl daemon-reload
#
# Даёт systemd больше времени на мягкую остановку после SIGTERM (до kill -9).

[Service]
TimeoutStopSec=120
